# sv/master/run
#!/bin/sh

sv start zookeeper || exit 1

exec 2>&1
exec mesos-master --zk=zk://localhost:2181/mesos --log_dir=/var/log --quiet=true
