# etc/supervisord-master.conf
[program:master]
command=mesos-master --zk=zk://localhost:2181/mesos
autorestart=true

