# etc/supervisord-zookeeper.conf
[program:zookeeper]
command=/zookeeper/bin/zkServer.sh start-foreground
autorestart=true

