# etc/supervisord-syslog.conf
[program:syslog]
command=/sbin/syslogd -n
autorestart=true

