# etc/supervisord-master.conf
[program:master]
command=mesos-master --zk=zk://localhost:2181/mesos --log_dir=/var/log --quiet=true
autorestart=true

