# etc/supervisord-chronos.conf
[program:chronos]
command=/chronos/bin/start-chronos.bash --master zk://localhost:2181/mesos --zk_hosts zk://localhost:2181/mesos --http_port 8081
autorestart=true

